# etc/systemd/system/backup.service
[Service]
ExecStart=/usr/local/sbin/backup
